# docs/ops/clamscan.timer
[Unit]
Description=Run nightly virus scan for document storage

[Timer]
OnCalendar=*-*-* 02:30:00
Persistent=true

[Install]
WantedBy=timers.target
